# infra/systemd/global-pulse-source-verify.timer
[Unit]
Description=Run Global Pulse source ingest verify hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true
Unit=global-pulse-source-verify.service

[Install]
WantedBy=timers.target
