# infra/systemd/global-pulse-snapshot.timer
[Unit]
Description=Run Global Pulse snapshot builder every 30 minutes

[Timer]
OnCalendar=*:0/30
RandomizedDelaySec=120
Persistent=true
Unit=global-pulse-snapshot.service

[Install]
WantedBy=timers.target
